# docker/supervisor-apt-cacher.conf
[program:apt-cacher-ng]
command=/usr/sbin/apt-cacher-ng ForeGround=1
